# scripts/backup/systemd/kkua-redis-backup.timer
[Unit]
Description=Run KKUA Redis backup daily  
Requires=kkua-redis-backup.service

[Timer]
# 매일 새벽 3시에 실행 (PostgreSQL 백업 후)
OnCalendar=*-*-* 03:00:00

# 시스템 부팅 후 20분 뒤에 한 번 실행 (누락된 백업 방지)
OnBootSec=20min

# 서비스가 실행되지 않았을 때만 실행
Persistent=true

# 무작위 지연 (서버 부하 분산)
RandomizedDelaySec=300

[Install]
WantedBy=timers.target