# scripts/backup/systemd/kkua-postgres-backup.timer
[Unit]
Description=Run KKUA PostgreSQL backup daily
Requires=kkua-postgres-backup.service

[Timer]
# 매일 새벽 2시에 실행
OnCalendar=*-*-* 02:00:00

# 시스템 부팅 후 15분 뒤에 한 번 실행 (누락된 백업 방지)
OnBootSec=15min

# 서비스가 실행되지 않았을 때만 실행
Persistent=true

# 무작위 지연 (서버 부하 분산)
RandomizedDelaySec=300

[Install]
WantedBy=timers.target