# scripts/backup/systemd/kkua-redis-backup.service
[Unit]
Description=KKUA Redis Database Backup
After=network.target
Wants=kkua-redis-backup.timer

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/kkua
ExecStart=/opt/kkua/scripts/backup_redis.sh

# 환경변수 파일 (필요시 수정)
EnvironmentFile=-/opt/kkua/backend/.env

# 로깅 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kkua-redis-backup

# 보안 설정
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/kkua/backups /opt/kkua/logs /tmp

# 실행 시간 제한 (15분)
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target